# root/etc/rc.local
#!/bin/sh

bash /root/mineros/bootstrap.sh &

exit 0
